# scripts/cron.d/moltis-telegram-webhook-monitor
# Moltis Telegram/Webhook Monitor Cron
# ====================================
# Runs continuous Telegram webhook + response cleanliness checks.
#
# Installed by: .github/workflows/deploy.yml (GitOps)
# Manual install: sudo cp /opt/moltinger-active/scripts/cron.d/moltis-telegram-webhook-monitor /etc/cron.d/

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Keep default monitor aligned with production polling mode.
# If webhook cutover is enabled intentionally, set TELEGRAM_REQUIRE_WEBHOOK=true explicitly.
TELEGRAM_REQUIRE_WEBHOOK=false
TELEGRAM_REQUIRE_HTTPS_WEBHOOK=true
TELEGRAM_REQUIRE_TEST_USER=false
TELEGRAM_PROBE_DISABLE_NOTIFICATION=true
TELEGRAM_MAX_PENDING_UPDATES=3
MONITOR_VALIDATE_PREVIEW=true
MOLTIS_BASE_URL=http://localhost:13131
MOLTIS_ENV_FILE=/opt/moltinger-active/.env

# Every 15 minutes run monitor and append JSON report.
# Active sendMessage probe runs only when TELEGRAM_TEST_USER is set explicitly.
*/15 * * * * root mkdir -p /var/log/moltis && /opt/moltinger-active/scripts/telegram-webhook-monitor.sh --json >> /var/log/moltis/telegram-webhook-monitor.log 2>&1
