# systemd/moltis-telegram-web-user-monitor.service
[Unit]
Description=Moltis Telegram Web user monitor probe
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/moltinger
EnvironmentFile=-/opt/moltinger/.env
Environment=MOLTIS_ENV_FILE=/opt/moltinger/.env
Environment=TELEGRAM_WEB_STATE=/opt/moltinger/data/.telegram-web-state.json
Environment=TELEGRAM_WEB_TARGET=@moltinger_bot
Environment=TELEGRAM_WEB_PROBE_PROFILE=echo_ping
Environment=TELEGRAM_WEB_MESSAGE=test2
Environment=TELEGRAM_WEB_TIMEOUT_SECONDS=45
Environment=TELEGRAM_WEB_MIN_REPLY_LEN=2
Environment=TELEGRAM_WEB_COMPOSER_RETRIES=2
ExecStart=/bin/bash -lc 'mkdir -p /var/log/moltis && /opt/moltinger/scripts/telegram-web-user-monitor.sh >> /var/log/moltis/telegram-web-user-monitor.log 2>&1'

[Install]
WantedBy=multi-user.target
