# systemd/moltis-telegram-web-user-monitor.timer
[Unit]
Description=Run Moltis Telegram Web user monitor every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=moltis-telegram-web-user-monitor.service

[Install]
WantedBy=timers.target
